# deploy/systemd/romad.service
[Unit]
Description=TagIt daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.tagit
ExecStart=%h/.local/bin/tagitd
Restart=on-failure
RestartSec=2
Environment=HOME=%h
Environment=PATH=%h/.local/bin:%h/.fnm/aliases/default/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=default.target
