# deploy/launchd/com.tagit.tagitd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tagit.tagitd</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>mkdir -p "$HOME/.tagit" &amp;&amp; cd "$HOME/.tagit" &amp;&amp; exec "$HOME/.local/bin/tagitd"</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>WorkingDirectory</key>
	<string>/Users/Shared</string>
</dict>
</plist>
